# crontabs
0 8,16 * * * node /home/node/app/index.js
# remember to end this file with an empty new line
